# scripts/laravel-queue-worker.service
[Unit]
Description=Laravel Queue Worker
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/turnos-app
ExecStart=/usr/bin/php artisan queue:work --queue=default --tries=3 --timeout=60 --memory=512 --sleep=3 --max-jobs=1000 --max-time=3600
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=laravel-queue-worker

# Configuración de recursos
LimitNOFILE=65536
MemoryMax=1G

# Variables de entorno (ajusta según tu configuración)
Environment=APP_ENV=production
Environment=APP_DEBUG=false

[Install]
WantedBy=multi-user.target
